# bmdns.service
[Unit]
Description: Bare-Metal DNS

[Service]
ExecStart=/usr/local/bin/bare-metal-dns

[Install]
WantedBy=default.target