# bmdns.service
[Unit]
Description=Bare-Metal DNS

[Service]
ExecStart=/usr/local/bin/bare-metal-dns
Restart=always

[Install]
WantedBy=default.target
